# deploy/supervisor/node.js.conf
[program:pingduino]
command=node /srv/pingduino/source/server/server.js
user=ubuntu
autostart=true
autorestart=true
redirect_stderr=True

